# backend-api-deb/DEBIAN/prerm
#!/bin/bash
systemctl stop backend-api.service
systemctl disable backend-api.service
exit 0
